# supervisord.conf
[supervisord]
nodaemon=true

[program:symfony]
command=symfony server:start --port=6969
autostart=true
autorestart=true

